# infra/vps/systemd/ping-watch-worker@.service
[Unit]
Description=Ping Watch worker (%i)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/ping-watch/%i/app/worker
EnvironmentFile=/etc/ping-watch/%i.env
ExecStart=/var/www/ping-watch/%i/app/worker/.venv/bin/python -m app.cli run --queue clip_uploaded
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
